# src/gunicorn.conf
description "Application server for flask-gunicorn"
start on runlevel [2345]
stop on runlevel [!2345]

respawn
setuid ubuntu
setgid ubuntu

script
    # ... environment settings
    # ...
    export PATH="/home/ubuntu/.pyenv/bin:$PATH"
    #echo $PATH
    #eval "$(pyenv init -)"
    #eval "$(pyenv virtualenv-init -)"
    cd /opt/verodin/src
    . /home/ubuntu/.pyenv/versions/verodin/bin/activate
    #echo `python --version`
    exec gunicorn cc:app  --bind 0.0.0.0:8317 --access-logfile /var/log/verodin/access.log --error-logfile /var/log/verodin/error.log
end script
