# src/poller.conf
description "AWS Poller"
start on runlevel [2345]
stop on runlevel [!2345]

respawn
setuid ubuntu
setgid ubuntu

script
    # ... environment settings
    # ...
    # Activate VirtEnv
    export PATH="/home/ubuntu/.pyenv/bin:$PATH"
    cd /opt/verodin/src
    . /home/ubuntu/.pyenv/versions/verodin/bin/activate

    exec python ccLoop.py
end script
